# scripts/com.codex.daily-financial-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.codex.daily-financial-brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>FINANCE_BRIEF_EMAIL="${FINANCE_BRIEF_EMAIL}" FINANCE_EMAIL_SENDER="${FINANCE_EMAIL_SENDER}" FINANCE_EMAIL_PASSWORD="${FINANCE_EMAIL_PASSWORD}" /Users/mofurong/.codex/skills/financial-news-collector/scripts/daily_brief.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>40</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/mofurong</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mofurong/Documents/每日财经新闻/daily_brief_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mofurong/Documents/每日财经新闻/daily_brief_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
